# deployment/argos-dev.service
[Unit]
Description=Argos Development Server
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/projects/ArgosFinal
ExecStart=/usr/bin/npm run dev
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=argos-dev
Environment="NODE_ENV=development"
Environment="PORT=5173"

[Install]
WantedBy=multi-user.target