# scripts/simple-keepalive.service
[Unit]
Description=Simple Development Server Keep-Alive
After=network.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/projects/ArgosFinal
ExecStart=/home/pi/projects/ArgosFinal/scripts/simple-keepalive.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=NODE_OPTIONS=--max-old-space-size=2048
Environment=PATH=/usr/local/bin:/usr/bin:/bin:/home/pi/.local/bin

[Install]
WantedBy=multi-user.target